# vizceral/rootfs/etc/service/node/run
#!/bin/sh


su www-data <<'EOF'
nep run dev
EOF
sleep 999d