# converter/rootfs/etc/service/agg/run
#!/bin/bash

/usr/bin/vrctlviz convert