# rootfs/etc/service/entrypoint/run
#!/bin/bash

exec /usr/bin/vrctlviz run