# vizceral/rootfs/etc/service/entrypoint/run
#!/bin/bash
source ~/.vrctlviz.yaml

npm run dev 