# vizceral/rootfs/etc/service/entrypoint/run
#!/bin/bash
source /.shared/status/vars

npm run dev 