# vizceral/rootfs/etc/service/agg/run
#!/bin/bash
source ~/.vrctlviz.yaml
/usr/bin/vrctlviz convert