# vizceral/rootfs/etc/service/agg/run
#!/bin/bash
# Source vars file if one exist
if [ -f ~/.vrctlvizcfg.yaml ]; then
source ~/.vrctlvizcfg.yaml
fi

/usr/bin/vrctlviz convert